# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info

[program:ai_service]
command=uvicorn main:app --host 0.0.0.0 --port 8000
directory=/app/ai_service
autostart=true
autorestart=true
startretries=5
stderr_logfile=/var/log/ai_service.err.log
stdout_logfile=/var/log/ai_service.out.log
environment=
    GOOGLE_API_KEY="%(ENV_GOOGLE_API_KEY)s",
    GEMINI_MODEL="%(ENV_GEMINI_MODEL)s",
    OPENWEATHER_KEY="%(ENV_OPENWEATHER_KEY)s",
    ORS_API_KEY="%(ENV_ORS_API_KEY)s"

[program:node_server]
command=node src/index.js
directory=/app/server
autostart=true
autorestart=true
startretries=5
startsecs=8
stderr_logfile=/var/log/node_server.err.log
stdout_logfile=/var/log/node_server.out.log
environment=
    PORT="5000",
    PYTHON_SERVICE_URL="http://localhost:8000",
    OPENWEATHER_KEY="%(ENV_OPENWEATHER_KEY)s",
    ORS_API_KEY="%(ENV_ORS_API_KEY)s",
    NODE_ENV="production"
